# backend/deploy/systemd/finme-api.service
[Unit]
Description=Finme Backend - API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=finme
Group=finme
WorkingDirectory=/opt/finme
ExecStart=/opt/finme/bin/finme-server api --config /opt/finme/config/config.toml
Restart=on-failure
RestartSec=3s
StartLimitInterval=300
StartLimitBurst=10

# 资源限制
LimitNOFILE=65535
MemoryMax=1G
TasksMax=512

# 安全收口
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/opt/finme/data

[Install]
WantedBy=multi-user.target
