# backend/deploy/systemd/finme-scheduler.service
[Unit]
Description=Finme Backend - Scheduler
After=finme-api.service

[Service]
Type=simple
User=finme
Group=finme
WorkingDirectory=/opt/finme
ExecStart=/opt/finme/bin/finme-server scheduler --config /opt/finme/config/config.toml
Restart=on-failure
RestartSec=3s

LimitNOFILE=8192
MemoryMax=512M
TasksMax=128

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/finme/data

[Install]
WantedBy=multi-user.target
